# dist/waarp-gatewayd.service
[Unit]
Description=Waarp Gateway server

[Service]
Type=simple
User=waarp
Group=waarp
WorkingDirectory=/var/lib/waarp-gateway
ExecStart=/bin/sh -c 'PATH=/usr/share/waarp-gateway:$PATH exec /usr/bin/waarp-gatewayd server -c /etc/waarp-gateway/gatewayd.ini'
Restart=on-failure

[Install]
WantedBy=multi-user.target

