# dist/waarp-gateway-get-remote.service
[Unit]
Description=Verifies files on sftp partners for Waarp Gateway

[Service]
Type=oneshot
User=waarp
Group=waarp
ExecStart=/usr/share/waarp-gateway/get-remotes
